# installation/ansible/roles/setup-daphne/files/daphne_restart.path
[Unit]
Description=Triggers the daphne service

[Path]
PathModified=/run/daphne/restart.txt

[Install]
WantedBy=multi-user.target
