# installation/ansible/roles/setup-daphne/files/daphne.socket
[Unit]
Description=daphne socket

[Socket]
ListenStream=/run/daphne/socket

[Install]
WantedBy=sockets.target
