# installation/ansible/roles/setup-gunicorn/files/gunicorn_restart.path
[Unit]
Description=Triggers the django3 gunicorn service

[Path]
PathModified=/run/gunicorn/gunicorn_restart.txt

[Install]
WantedBy=multi-user.target
