# backend-django/installation/ansible/roles/setup-gunicorn/files/gunicorn_restart.service
[Unit]
Description=Restart django3 gunicorn service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart gunicorn
