# backend-django/installation/ansible/roles/setup-gunicorn/files/gunicorn.socket
[Unit]
Description=django3 gunicorn socket

[Socket]
ListenStream=/run/gunicorn/socket

[Install]
WantedBy=sockets.target
