# backend-django/installation/ansible/roles/setup-daphne/files/daphne_restart.service
[Unit]
Description=Restart daphne service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart daphne
